# deploy/supervisor/mirrors.ini
[program:mirrors]
command=/data/prd/mirrors/env/bin/python src/run.py
directory=/data/prd/mirrors
user=daimon
autostart=true
autorestart=true
startretries=5
stdout_logfile=/data/prd/mirrors/logs/supervisor_stdout.log
stderr_logfile=/data/prd/mirrors/logs/supervisor_stderr.log
environment=DEBUG="False"
